# install/autobot.service
[Unit]
Description=SRE Autobot
 
[Service]
ExecStart=/usr/bin/python3 __path__
 
[Install]
WantedBy=multi-user.target
 
 
 
