# install/autobot.service
[Unit]
Description=SRE Autobot
 
[Service]
ExecStart=__path__
 
[Install]
WantedBy=multi-user.target
 
 
 
