# install/sre.service
[Unit]
Description=SRE Autobot
 
[Service]
ExecStart=/usr/bin/python3 __path__ --daemon
 
[Install]
WantedBy=multi-user.target
 
 
 
