# srebot/datafiles/sre.service
[Unit]
Description=SRE Bot __name__

[Service]
Restart=Always
ExecStart=__exec_start__

[Install]
WantedBy=multi-user.target
